# vegeta/supervisord.conf
[supervisord]
nodaemon=true

[program:vegeta]
command=/bin/vegeta-runner.sh
